# offregister_taiga/data/circusd.conf
description "circusd"

start on filesystem and net-device-up IFACE=lo
stop on runlevel [016]

respawn
exec "%(CIRCUS_VENV)s/bin/circusd" "%(CONF_DIR)s/circus.ini"
